# systemd/hawkdog.service
[Unit]
Description=hawkdog (file access alert daemon)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/hawkdog
ExecStart=%h/.local/share/hawkdog/hawkdog
Restart=always
RestartSec=2

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/hawkdog %h/.config/sentinel-watch %h/.local/share/hawkdog %h/.clawdbot/credentials

[Install]
WantedBy=default.target
